# apps/autotest/ub_autotest.socket
[Socket]
ListenStream=127.0.0.1:8882
NoDelay=true
KeepAlive=true
 
[Install]
WantedBy=sockets.target